# systemd/llama-cpp-server.service
[Unit]
Description=llama.cpp Inference Server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/Data/Programs/AI/llama-cpp-server
ExecStart=/usr/bin/python3 %h/Data/Programs/AI/llama-cpp-server/launcher.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
